# etc/systemd/system/pq-save.timer
[Unit]
Description=Run quota save every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target

