# etc/systemd/system/vless-gc.service
[Unit]
Description=GC expired temporary VLESS nodes
After=local-fs.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vrctl vless-gc

